# crontab
* * * * * cd /var/www/localhost/htdocs && php artisan schedule:run >> /proc/self/fd/1 2>/proc/self/fd/2